# deploy/linux/systemd/japan-job-agent-backup.service
[Unit]
Description=Japan Job Agent PostgreSQL backup
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/japan-job-agent/tmp/current-release
ExecStart=/usr/bin/docker compose --env-file /opt/japan-job-agent/deploy/linux/.env -f deploy/linux/compose.yml --profile maintenance run --rm backup
