# deploy/linux/systemd/japan-job-agent-public-ats-refresh.service
[Unit]
Description=Refresh persisted public ATS Discovery tenants
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/japan-job-agent
TimeoutStartSec=2h
ExecStart=/usr/bin/docker compose --env-file deploy/linux/.env -f deploy/linux/compose.yml run --rm worker pnpm discovery:public-ats
