# deploy/linux/systemd/japan-job-agent-backup.timer
[Unit]
Description=Daily Japan Job Agent PostgreSQL backup timer

[Timer]
OnCalendar=*-*-* 03:20:00 Asia/Tokyo
Persistent=true
RandomizedDelaySec=900

[Install]
WantedBy=timers.target
