# deploy/linux/systemd/japan-job-agent-freshness.timer
[Unit]
Description=Daily Japan Job Agent publication freshness timer

[Timer]
OnCalendar=*-*-* 08:30:00 Asia/Tokyo
Persistent=true
RandomizedDelaySec=900

[Install]
WantedBy=timers.target
