# deploy/linux/systemd/japan-job-agent-promotion.service
[Unit]
Description=Resolve and promote verified Discovery candidates
Requires=docker.service
After=docker.service japan-job-agent-discovery-refresh.service japan-job-agent-public-ats-refresh.service

[Service]
Type=oneshot
WorkingDirectory=/opt/japan-job-agent
TimeoutStartSec=6h
ExecStart=/usr/bin/docker compose --env-file deploy/linux/.env -f deploy/linux/compose.yml run --rm worker pnpm discovery:promote-talentio
ExecStart=/usr/bin/docker compose --env-file deploy/linux/.env -f deploy/linux/compose.yml run --rm worker pnpm discovery:promote-public-ats
ExecStart=/usr/bin/docker compose --env-file deploy/linux/.env -f deploy/linux/compose.yml run --rm worker pnpm temporal:schedules
