# deploy/linux/systemd/japan-job-agent-discovery-refresh.timer
[Unit]
Description=Daily Japan Job Agent Discovery refresh timer

[Timer]
OnCalendar=*-*-* 04:10:00 Asia/Tokyo
Persistent=true
RandomizedDelaySec=900

[Install]
WantedBy=timers.target
