# deploy/linux/systemd/japan-job-agent-public-ats-refresh.timer
[Unit]
Description=Daily public ATS Discovery refresh timer

[Timer]
OnCalendar=*-*-* 05:50:00 Asia/Tokyo
Persistent=true
RandomizedDelaySec=900

[Install]
WantedBy=timers.target
